# docker/cron/crontab
* * * * * root cd /app/bin && /usr/local/bin/php console GenerateSitemap 2> /var/log/cron.log
* * * * * root date >> /var/log/test-cron.log
